# pybob/system_config/upstart/cb_logwatcher.conf
#!upstart

description "Chez Bob log watcher"

start on startup
stop on shutdown

setuid sodasrv
setgid sodasrv

chdir /git/pybob/logwatcher
exec python3.4 logwatcher.py

