# pybob/system_config/upstart/cb_transaction_watcher.conf
#!upstart

description "Chez Bob transaction watcher"

start on startup
stop on shutdown

setuid sodasrv
setgid sodasrv

exec /git/pybob/notifiers/transaction_watcher

